# services/srcs/mySql/srcs/supervisor.conf
[supervisord]
nodaemon=true


[program:mysql]
command = sh -c "mysqld -uroot && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "/telegraf-1.17.2/usr/bin/telegraf --config /telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"
# [program:sshd]
# command = sh -c "/usr/sbin/sshd -D && kill -s SIGTERM $(cat supervisord.pid)"

# [program:php-fpm7]
# command = sh -c "php-fpm7 -F && kill -s SIGTERM $(cat supervisord.pid)"

# [program:nginx]
# command = sh -c "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"