# services/srcs/influxdb/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:influx]
command = sh -c "influxd && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "/telegraf-1.17.2/usr/bin/telegraf --config /telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"
